# config/systemd/user/cros-garcon.service.d/override.conf
[Service]
Environment="PATH=%h/.local/bin:%h/.local/share/cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/local/games:/usr/sbin:/usr/bin:/usr/games:/sbin:/bin"
